# deploy/compass.service
[Unit]
Description=AABW Compass (FastAPI + Bedrock agent)
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/aabw-compass
Environment=AWS_REGION=us-east-1
Environment=BEDROCK_MODEL_ID=us.anthropic.claude-sonnet-4-5-20250929-v1:0
ExecStart=/home/ec2-user/aabw-compass/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
